# sched-ext-userland.service
[Unit]
Description=Start scx_userland

[Service]
Type=oneshot
ExecStart=/usr/bin/scx_userland

[Install]
WantedBy=multi-user.target
