# sched-ext-pair.service
[Unit]
Description=Start scx_pair

[Service]
Type=oneshot
ExecStart=/usr/bin/scx_pair

[Install]
WantedBy=multi-user.target
