# sched-ext-flatcg.service
[Unit]
Description=Start scx_flatcg

[Service]
Type=oneshot
ExecStart=/usr/bin/scx_flatcg

[Install]
WantedBy=multi-user.target
