# sched-ext-rustland.service
[Unit]
Description=Start scx_rustland

[Service]
Type=oneshot
ExecStart=/usr/bin/scx_rustland

[Install]
WantedBy=multi-user.target
