# systemd-services/sched-ext-userland.service
[Unit]
Description=Start scx_userland
ConditionPathExists=/sys/kernel/debug/sched/ext

[Service]
Type=simple
ExecStart=/bin/sh -c '/usr/bin/scx_userland 2>&1 | grep -E "warn|error"'
Restart=always

[Install]
WantedBy=multi-user.target
