# systemd-services/sched-ext-rusty.service
[Unit]
Description=Start scx_rusty
ConditionPathExists=/sys/kernel/debug/sched/ext

[Service]
Type=simple
ExecStart=/bin/sh -c '/usr/bin/scx_rusty 2>&1 | grep -E "warn|error"'
Restart=always

[Install]
WantedBy=multi-user.target
